# scripts/lightning-monitor.service
[Unit]
Description=CoinPay Lightning Payment Monitor
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/src/coinpayportal
ExecStart=/home/ubuntu/src/coinpayportal/.venv/bin/python3 scripts/lightning-monitor.py
EnvironmentFile=/home/ubuntu/src/coinpayportal/.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
